# services/kakao/deploy/kittykakao.service
[Unit]
Description=KittyKakao — KakaoTalk WebSocket relay
After=network.target

[Service]
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{REMOTE_DIR}}
EnvironmentFile={{REMOTE_DIR}}/.env
ExecStart={{REMOTE_DIR}}/kittykakao
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
